# Server/deploy/gunicorn_super.ini
[program:apppython]
directory=/blog/Server
command=gunicorn -c ./deploy/gunicorn_config.py app.main:app
startsecs=3
stdout_logfile=/blog/Server/log/supervisor.log
